# supervisor/cloudwatch.conf
[program:cloudwatch]
command=/var/awslogs/bin/awslogs-agent-launcher.sh
autostart=true
autorestart=true
;stderr_logfile=/var/log/error.log
;stdout_logfile=/var/log/access.log