# supervisor/worker.conf
[program:worker]
directory=/path/to/project/directory/
command=/usr/local/bin/celery -A fxtrain worker --pool=solo
autostart=true
autorestart=true
;stderr_logfile=/var/log/error.log